# flakes/ansible/server/files/services/sshd-alert.service
[Unit]
Description=sshd-alert service

[Service]
ExecStart=/opt/scripts/sshd-alert.sh

[Install]
WantedBy=multi-user.target
